# config/autostart/elfin-browser.desktop
[Desktop Entry]
Type=Application
Name=Elfin Browser
Exec=elfin-browser
X-GNOME-Autostart-enabled=true
Hidden=false
NoDisplay=false
